# seleniumNodes/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:seleniumNode]
command=./opt/startNode.sh
